# omnivore/containers/systemd/omnivore-api.container
[Unit]
Description=Omnivore API Server
PartOf=omnivore.target
Requires=postgres.service garage.service omnivore-redis.service omnivore-migrate.service
After=postgres.service garage.service omnivore-redis.service omnivore-migrate.service

{{#if (command_success "uname -r | grep -qi wsl")}}
[Quadlet]
DefaultDependencies=false

{{/if}}
[Service]
Restart=always

[Container]
Image=ghcr.io/omnivore-app/sh-backend:latest

Secret=omnivore-app-password,type=env,target=PG_PASSWORD
Secret=omnivore-jwt-secret,type=env,target=JWT_SECRET
Secret=omnivore-sso-jwt-secret,type=env,target=SSO_JWT_SECRET
Secret=omnivore-s3-access-key,type=env,target=AWS_ACCESS_KEY_ID
Secret=omnivore-s3-secret-key,type=env,target=AWS_SECRET_ACCESS_KEY

Network=omnivore.network
Network=traefik.network
NetworkAlias=omnivore-api

Environment=API_ENV=local
Environment=PG_HOST=postgres
Environment=PG_USER=app_user
Environment=PG_DB=omnivore
Environment=PG_PORT=5432
Environment=PG_POOL_MAX=20

# URLs - adjust for your domain
Environment=CLIENT_URL=https://omnivore.{{domain}}
Environment=GATEWAY_URL=http://omnivore-api:8080/api
Environment=CONTENT_FETCH_URL=http://content-fetch:8080/?token=omnivore_fetch_token
Environment=IMAGE_PROXY_URL=https://omnivore-api.{{domain}}/images

# Redis
Environment=REDIS_URL=redis://redis:6379/0

# S3/Garage
Environment=GCS_USE_LOCAL_HOST=true
Environment=GCS_UPLOAD_BUCKET=omnivore
Environment=AWS_REGION=garage
Environment=AWS_S3_ENDPOINT_URL=http://garage:3900
Environment=LOCAL_MINIO_URL=http://garage:3900

# Features
Environment=AUTO_VERIFY=true
Environment=CONTENT_FETCH_QUEUE_ENABLED=true

HealthCmd=nc -z 0.0.0.0 8080
HealthInterval=15s
HealthTimeout=10s
HealthRetries=6

# Traefik labels
Label=traefik.enable=true
Label=traefik.docker.network=systemd-traefik

# HTTP -> HTTPS redirect
Label=traefik.http.routers.omnivore-api-http.entrypoints=http
Label=traefik.http.routers.omnivore-api-http.rule=Host(`omnivore-api.{{domain}}`)
Label=traefik.http.routers.omnivore-api-http.middlewares=redir-https@file
Label=traefik.http.routers.omnivore-api-http.service=noop@internal

# HTTPS route
Label=traefik.http.routers.omnivore-api-https.entrypoints=https
Label=traefik.http.routers.omnivore-api-https.rule=Host(`omnivore-api.{{domain}}`)
Label=traefik.http.routers.omnivore-api-https.tls=true
Label=traefik.http.routers.omnivore-api-https.middlewares=gzip@file
Label=traefik.http.services.omnivore-api.loadbalancer.server.port=8080

[Install]
WantedBy=omnivore.target
