# superset/systemd/user/superset.target
[Unit]
Description=Superset Application Stack
