# langfuse/systemd/user/langfuse.target
[Unit]
Description=Langfuse Application Stack
