# plane/systemd/user/plane.target
[Unit]
Description=Plane Application Stack

[Install]
WantedBy=default.target
