# traefik/systemd/user/http.socket
[Socket]
ListenStream=80
FileDescriptorName=http
Service=traefik.service

[Install]
WantedBy=sockets.target
