# omnivore/systemd/user/omnivore.target
[Unit]
Description=Omnivore Read-it-Later Service Stack

[Install]
WantedBy=default.target
