# plane/containers/systemd/plane-space.container
[Unit]
Description=Plane Space (Public Sharing)
PartOf=plane.target
After=plane-api.service plane-worker.service plane-web.service

{{#if (command_success "uname -r | grep -qi wsl")}}
[Quadlet]
DefaultDependencies=false

{{/if}}
[Service]
Restart=always

[Container]
Image=artifacts.plane.so/makeplane/plane-space:stable
Network=plane.network
NetworkAlias=space

# Dozzle group
Label=dev.dozzle.group=quadlet

[Install]
WantedBy=plane.target
